# backup-source/envoy-ingress-with-jwt/supervisord.conf
[program:envoy_blue]
process_name=envoy_blue
command=/etransfer/esg/envoy.sh
autostart=true
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_events_enabled=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_events_enabled=true

[program:envoy_green]
process_name=envoy_green
command=/etransfer/esg/envoy.sh
autostart=false
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_events_enabled=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stderr_events_enabled=true

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0

[unix_http_server]
file=%(here)s/supervisor.sock

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# [eventlistener:envoy_blue_fatal]
# process_name=envoy_blue
# command=kill -SIGQUIT $PPID
# events=PROCESS_STATE_FATAL

# [eventlistener:envoy_green_fatal]
# process_name=envoy_green
# command=kill -SIGQUIT $PPID
# events=PROCESS_STATE_FATAL
